# gnome/config/autostart/xbindkeys.desktop
[Desktop Entry]
Name=xbindkeys
Exec=xbindkeys
Terminal=false
Type=Application
X-Gnome_autostart-enabled=true
